# rpi_sound_trigger/install/sound-trigger.service
[Unit]
Description=RPi Sound Trigger (mic → MQTT display alert)
After=network-online.target mosquitto.service sound.target graphical.target
Wants=network-online.target mosquitto.service

[Service]
Type=simple
User=pi
Group=pi
WorkingDirectory=/home/pi/ESP32-8048S070_movie_player/rpi_sound_trigger
Environment=PYTHONUNBUFFERED=1
# Desktop session access for the on-screen VU meter (ignored when headless / no HDMI).
Environment=DISPLAY=:0
Environment=XAUTHORITY=/home/pi/.Xauthority
Environment=XDG_RUNTIME_DIR=/run/user/1000
# Auto UI when a screen is connected; headless otherwise (override with --no-display / --display).
ExecStart=/home/pi/ESP32-8048S070_movie_player/rpi_sound_trigger/.venv/bin/python \
  /home/pi/ESP32-8048S070_movie_player/rpi_sound_trigger/src/main.py \
  -c /home/pi/ESP32-8048S070_movie_player/rpi_sound_trigger/config.yaml
Restart=on-failure
RestartSec=3

# Mic/webcam + DRM/framebuffer for the VU meter
SupplementaryGroups=audio video render

[Install]
WantedBy=multi-user.target
